# .config/systemd/user/desktop-bg.service
[Unit]
Description=Arts & Culture Desktop Service

[Service]
Type=oneshot
ExecStart=/home/julian/.local/bin/artsandculturedesktop
