# .config/systemd/user/maildir-archive.timer
[Unit]
Description=Yearly maildir archive (runs Jan 15)

[Timer]
OnCalendar=*-01-15
Persistent=true

[Install]
WantedBy=timers.target
