# .config/systemd/user/maildir-archive.service
[Unit]
Description=Maildir archive by year

[Service]
Type=oneshot
ExecStart=%h/.local/bin/maildir-archive
