# .config/systemd/user/mbsync.timer
[Unit]
Description=Mail sync timer

[Timer]
OnUnitActiveSec=2m
Unit=mbsync@fastmail.service

[Install]
WantedBy=timers.target
