# .config/systemd/user/notify-failure@.service
[Unit]
Description=OnFailure email for %i

[Service]
Type=oneshot
ExecStart=/usr/bin/mail -s "Service %i failed" julian
