# .config/systemd/user/notmuch.timer
[Unit]
Description=Periodic notmuch index update

[Timer]
OnCalendar=minutely
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target
