# .config/systemd/user/selenite-update.service
[Unit]
Description=Selenite Lamp Update Service

[Service]
Type=oneshot
ExecStart=/home/julian/.local/bin/update-lamp-state
