# .config/systemd/user/selenite-update.timer
[Unit]
Description=Timer for the Selenite service

[Timer]
OnActiveSec=0
OnUnitActiveSec=10
AccuracySec=5

[Install]
WantedBy=dev-selenite\x2dlamp.device
