# .config/systemd/user/selenite-update.service
[Unit]
Description=Selenite Lamp Update Service

[Service]
ExecStart=/usr/bin/selenite-lamp-updater -c %E/selenite-lamp-updater/config.toml
Restart=always

[Install]
WantedBy=dev-selenite\x2dlamp.device
