# .config/systemd/user/selenite.service
[Unit]
Description=Selenite Lamp Service

[Service]
Type=oneshot
ExecStart=/home/julian/.local/bin/update-lamp-state
