# .config/systemd/user/site-status.service
[Unit]
Description=Site Status Service

[Service]
Type=oneshot
ExecStart=/home/julian/.local/bin/update-site-status jellyfin.seemyvest.net
