# .config/systemd/user/tarsnap.service
[Unit]
Description=Run tarsnap backups
OnFailure=notify-failure@%n

[Service]
Type=oneshot
ExecStart=/home/julian/.local/bin/tarsnap-backup
