# .config/systemd/user/update-sbuild-chroots.service
[Unit]
Description=Update sbuild chroot tarballs

[Service]
Type=oneshot
ExecStart=%h/.local/bin/update-sbuild-chroots
