# .config/systemd/user/vdirsyncer.service.d/override.conf
[Service]
LoadCredentialEncrypted=fastmail-app-password:/home/julian/.config/credstore/fastmail-app-password.cred
