# .dotfiles-notes/udev-rules/70-monitor.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/home/julian/.local/bin/set-monitors"
